# service/docker/run
#!/bin/sh
exec docker daemon $DOCKER_OPTS
