# service/docker/run
#!/bin/sh

# If we're started with a DOCKER_HOST variable, don't start dockerd
if [ "${DOCKER:-}" = "no" ]; then
	sv stop "$PWD"
	exit
fi
exec docker daemon $DOCKER_OPTS
