# service/jenkins/run
#!/bin/sh
exec su -c /bin/jenkins jenkins
